# app/deploy/hydrohalo.service
[Unit]
Description=HydroHalo Swim Trainer
After=graphical.target network.target
Wants=graphical.target

[Service]
User=pi
WorkingDirectory=/home/pi/hydrohalo
ExecStart=/usr/bin/python3 -m app.main
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
